# packages/scingestor/scingestor-0.20.0.tar.gz/scingestor-0.20.0/config/scingestor.service
[Unit]
Description=SciCat Dataset Ingestor
After=network.target

[Service]
Type=simple
User=root
Nice=10
WorkingDirectory=/tmp/
ExecStart=/usr/bin/python3 /usr/bin/scicat_dataset_ingestor -c /etc/scingestor/default.yaml
Restart=always

[Install]
WantedBy=multi-user.target
